# auto-update.service
[Unit]
Description=WWT Music Club Auto-Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/wwtmc

# Set one of these:
# Option 1: Monitor a git repository
# Environment="SOURCE_REPO=https://github.com/user/wwtmc.git"
# Environment="CHECK_INTERVAL=3600"

# Option 2: Monitor a local directory
Environment="SOURCE_DIR=/mnt/updates"
Environment="CHECK_INTERVAL=900"

ExecStart=/bin/bash /var/www/wwtmc/auto-update.sh

# Restart policy
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wwtmc-auto-update

# Security
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
