# auto-update.service
[Unit]
Description=WWT Music Club Auto-Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/wwtmc

Environment="SOURCE_REPO=https://github.com/cn2450-netizen/wwtm-boosters-webpage.git"
Environment="CHECK_INTERVAL=900"

ExecStart=/bin/bash /var/www/wwtmc/auto-update.sh

# Restart policy
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wwtmc-auto-update

# Security
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
